# crontab
0 * * * * python3 /src/edward.py -t twitter
30 * * * * python3 /src/edward.py -t reddit
