# crontab
0 0 * * * python3 /src/edward.py -t reddit
0 * * * * python3 /src/edward.py -t twitter
