# crontab
0 * * * * python3 /src/edward.py -t reddit
*/10 * * * * python3 /src/edward.py -t twitter
