# crontab
0 */2 * * * python3 /src/edward.py -t twitter
0 */6 * * * python3 /src/edward.py -t reddit
